# docker/cockpit/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
user=root

[program:dbus]
command=/usr/bin/dbus-daemon --system --nofork --nopidfile
priority=10
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/dbus.log
stderr_logfile=/var/log/supervisor/dbus-error.log

[program:sshd]
command=/usr/sbin/sshd -D -e
priority=15
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sshd.log
stderr_logfile=/var/log/supervisor/sshd-error.log

[program:pmcd]
command=/usr/libexec/pcp/lib/pmcd start-foreground
priority=17
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/pmcd.log
stderr_logfile=/var/log/supervisor/pmcd-error.log

[program:pmlogger]
command=/usr/libexec/pcp/bin/pmlogger -t 60sec -c config.default -l /var/log/pcp/pmlogger/pmlogger.log /var/log/pcp/pmlogger/hostname
priority=18
autostart=true
autorestart=true
user=pcp
stdout_logfile=/var/log/supervisor/pmlogger.log
stderr_logfile=/var/log/supervisor/pmlogger-error.log

[program:cockpit-ws]
command=/usr/lib/cockpit/cockpit-ws --no-tls --port 9090 --local-ssh
priority=20
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/supervisor/cockpit-ws.log
stderr_logfile=/var/log/supervisor/cockpit-ws-error.log

[program:portainer]
command=/portainer --data /portainer_data --bind :9000 --hide-label owner=websoft9
priority=25
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/supervisor/portainer.log
stderr_logfile=/var/log/supervisor/portainer-error.log

[program:gitea]
command=/usr/local/bin/gitea web --config /etc/gitea/app.ini
directory=/var/lib/gitea
environment=HOME="/var/lib/gitea",USER="git",GITEA_WORK_DIR="/var/lib/gitea"
priority=26
autostart=true
autorestart=true
user=git
stdout_logfile=/var/log/supervisor/gitea.log
stderr_logfile=/var/log/supervisor/gitea-error.log

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
priority=30
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-error.log

[program:apphub]
command=python3 -m uvicorn src.main:app --host 0.0.0.0 --port 8080
directory=/websoft9/apphub
priority=40
autostart=true
autorestart=true
user=%(ENV_ADMIN_USER)s
stdout_logfile=/var/log/supervisor/apphub.log
stderr_logfile=/var/log/supervisor/apphub-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:apphub-cron]
command=cron -f
priority=42
autostart=true
autorestart=true
user=%(ENV_ADMIN_USER)s
stdout_logfile=/var/log/supervisor/apphub-cron.log
stderr_logfile=/var/log/supervisor/apphub-cron-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
