# docker/proxy/s6/init_user/run
#!/bin/bash

exec "/etc/s6-overlay/s6-rc.d/init_user/setup" "$@"