# docker/proxy/s6/init_user/run
#!/bin/bash

bash /etc/s6-overlay/s6-rc.d/init_user/init_user.sh

exec s6-svc -D /etc/s6/user