# docker/proxy/s6/init_user/up
exec s6-svc -D /etc/s6-overlay/s6-rc.d/init_user