# docker/s6/websoft9-supervisor/run
#!/command/with-contenv bash

set -euo pipefail

exec /websoft9/script/platform-entrypoint.sh