# build/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid
user=root
# 优雅退出配置
stopwaitsecs=30
# 最小可用文件描述符
minfds=1024

; XML-RPC server for services plugin (Story 6.1)
; Note: For development, using simple password. In production, use ENV variable.
[inet_http_server]
port=127.0.0.1:9001
username=admin
password=%(ENV_SUPERVISOR_PASSWORD)s

[supervisorctl]
serverurl=http://127.0.0.1:9001
username=admin
password=%(ENV_SUPERVISOR_PASSWORD)s

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ======================
; 系统基础服务
; ======================

[program:dbus]
command=/usr/bin/dbus-daemon --system --nofork --nopidfile
priority=10
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/dbus.log
stderr_logfile=/var/log/supervisor/dbus-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

[program:sshd]
command=/usr/sbin/sshd -D -e
priority=15
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/sshd.log
stderr_logfile=/var/log/supervisor/sshd-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

; ======================
; 业务服务
; ======================

[program:cockpit-ws]
command=/usr/lib/cockpit/cockpit-ws --no-tls --port 9090 --local-ssh
priority=20
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=15
user=root
stdout_logfile=/var/log/supervisor/cockpit-ws.log
stderr_logfile=/var/log/supervisor/cockpit-ws-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

[program:portainer]
command=/portainer --data /websoft9/data/portainer --bind :9000 --hide-label owner=websoft9
priority=25
autostart=true
autorestart=true
startsecs=15
startretries=3
stopwaitsecs=20
user=root
stdout_logfile=/var/log/supervisor/portainer.log
stderr_logfile=/var/log/supervisor/portainer-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

[program:gitea]
command=/usr/local/bin/gitea web --config /etc/gitea/app.ini
directory=/websoft9/data/gitea
environment=HOME="/websoft9/data/gitea",USER="git",GITEA_WORK_DIR="/websoft9/data/gitea"
priority=26
autostart=true
autorestart=true
startsecs=15
startretries=3
stopwaitsecs=20
user=git
stdout_logfile=/var/log/supervisor/gitea.log
stderr_logfile=/var/log/supervisor/gitea-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
priority=30
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopsignal=QUIT
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

[program:convex]
command=/usr/local/bin/convex-backend
directory=/websoft9/data/baas
environment=CONVEX_CLOUD_ORIGIN="http://127.0.0.1:3210",CONVEX_SITE_ORIGIN="http://127.0.0.1:3211",RUST_LOG="info",INSTANCE_NAME="websoft9"
priority=40
autostart=true
autorestart=true
startsecs=15
startretries=5
stopwaitsecs=20
user=root
stdout_logfile=/var/log/supervisor/convex.log
stderr_logfile=/var/log/supervisor/convex-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
