# docker/proxy/config/supervisord.conf
[supervisord]
nodaemon=false

[program:nginxproxymanage]
command=/init
autostart=true
directory=/