# docker/git/src/s6/user/finish
#!/bin/bash

s6-svscanctl -t /etc/s6/user
exit 100
