# docker/git/src/s6/user/finish
#!/bin/bash

cred_path="/data/gitea/credential"
if [ -e "$cred_path" ]; then
  echo "File $cred_path exists. disable /etc/s6/user"
  exit 125
else
  exit 0
fi
