# docker/git/src/s6/user/run
#!/bin/bash

[[ -f ./setup ]] && source ./setup
pushd /root > /dev/null
popd