# docker/git/src/s6/user/run
#!/bin/execlineb -P
s6-oneshot-runner

[[ -f ./setup ]] && source ./setup
pushd /root >/dev/null
exec s6-svc -D /etc/s6/user
popd