# docker/git/src/s6/user/run
#!/bin/bash

[[ -f ./setup ]] && source ./setup
exec s6-svc -D /etc/s6/user
popd
